# docker/conf/defaults/crontab
# speedtest crons
* * * * * php /config/www/artisan schedule:run >> /config/log/homelab/cron.log
*/5 * * * * php /config/www/artisan queue:retry all >> /config/log/homelab/cron.log